# etc/local.d/boot.start
sudo -u king pulseaudio --start -D&

/etc/local.d/wifi &
